# scripts/ai.hermes.secure-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.hermes.secure-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__INSTALL_DIR__/scripts/cron-update.sh</string>
    </array>

    <!-- Run daily at 04:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HERMES_HOME</key>
        <string>__HERMES_HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HERMES_HOME__/secure-update/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__HERMES_HOME__/secure-update/launchd.error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
